# crontab
* * * * * python3 /custom_apps/zeus_crm_api/process_file_uploads.py

# An empty line is required at the end of this file for a valid cron file.
